# .config/autostart/kdeconnectd_1_.desktop
[Desktop Entry]
Exec=/usr/lib/kdeconnectd
Icon=
Name=kdeconnectd
Path=
Terminal=False
Type=Application
