# .config/autostart/kdeconnectd.desktop
[Desktop Entry]
Exec=/usr/lib/kdeconnectd
Icon=dialog-scripts
Name=kdeconnectd
Path=
Type=Application
X-KDE-AutostartScript=true
